# .config/autostart/glipper.desktop
[Desktop Entry]
Comment=A clipboard history manager
Exec=glipper
Name=Glipper
OnlyShowIn=GNOME;Unity;
Type=Application
Hidden=true
