# .config/autostart/glipper.desktop
[Desktop Entry]
Name=Glipper
Comment=A clipboard history manager
Exec=glipper
Type=Application
OnlyShowIn=GNOME;Unity;
Hidden=true
